# crontab
*/10 * * * * python /app/etl.py >> /var/log/cron.log 2>&1

